# scs/deploy/nais-scs.service
[Unit]
Description=NAIS Secure Channel Services (SCS)
Documentation=https://github.com/your-org/nais-client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nais-scs
Group=nais-scs
ExecStart=/usr/local/bin/nais-scs --config /etc/nais-scs/config.toml
Restart=always
RestartSec=5
TimeoutStopSec=10

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes

# Allow write to data directory
ReadWritePaths=/var/lib/nais-scs

# Working directory
WorkingDirectory=/var/lib/nais-scs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nais-scs

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
